# systemd/mynodeone-sync-controller.service
[Unit]
Description=MyNodeOne Sync Controller - Event-Driven Config Push
After=network.target k3s.service
Requires=k3s.service

[Service]
Type=simple
User=root
WorkingDirectory=/path/to/MyNodeOne
ExecStart=/path/to/MyNodeOne/scripts/lib/sync-controller.sh watch
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mynodeone-sync

# Security
PrivateTmp=yes
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
